# hw2/p4/hw2_p4.timer
[Unit]
Description=Run hw2_p4 service every 10 minutes

[Timer]
OnCalendar=*:0/10
AccuracySec=1ms
